# systemd/user/gpg-agent.service
[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon --enable-ssh-support

[Install]
WantedBy=default.target
